# local-stt.service
[Unit]
Description=Local GPU Speech-to-Text Daemon
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/whisperubuntu-daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
